# packaging/udev/50-rcraid.rules.in
# SPDX-License-Identifier: GPL-2.0-only
#
# rcraid auto-bind rule — generated by install-dkms.sh
#
# When the kernel discovers a PCI device matching AMD's NVMe RAID Bottom
# (vendor 0x1022, device 0xb000) AND a subsystem_vendor that we recognise
# as an array member, hand it to /usr/sbin/rcraid-bind which transfers
# ownership from `nvme` to `rcbottom`.
#
# IMPORTANT: AMD's chipset shadows EVERY downstream NVMe behind 0x1022:0xb000,
# so the subsystem_vendor match below is what prevents this rule from
# grabbing the OS drive.  Edit @SUBSYSTEM_VENDOR@ if you change array
# hardware.

ACTION=="add|change", SUBSYSTEM=="pci", \
  ATTR{vendor}=="0x1022", ATTR{device}=="0xb000", \
  ATTR{subsystem_vendor}=="@SUBSYSTEM_VENDOR@", \
  RUN+="/usr/sbin/rcraid-bind %k"
